# templates/documents/openquad.container
[Unit]
Description=OpenQuad Documents Agent

[Container]
ContainerName=openquad-documents
Image=ghcr.io/myos-dev/openquad-documents:latest
Pull=missing
Network=openquad.network

PublishPort=127.0.0.1:18789:18789
UserNS=keep-id
NoNewPrivileges=true
DropCapability=ALL
ReadOnly=true
Tmpfs=/tmp:size=512M,mode=1777
WorkingDir=/home/node
Environment=HOME=/home/node
Environment=OPENQUAD_TEMPLATE=documents
Environment=OPENQUAD_IMAGE_NAME=openquad-documents
Environment=OPENCLAW_STATE_DIR=/home/node/.openclaw
Environment=OPENCLAW_CONFIG_PATH=/home/node/.openclaw/openclaw.json
Environment=OPENQUAD_DEFAULT_CONFIG=/usr/share/openquad/templates/documents/openclaw.json5
Environment=OLLAMA_BASE_URL=http://ollama:11434
Environment=DOCUMENT_WORKSPACE=/home/node/.openclaw/workspace/documents
Environment=NPM_CONFIG_CACHE=/home/node/.openclaw/.npm
Environment=XDG_CONFIG_HOME=/home/node/.openclaw/.config
Environment=XDG_DATA_HOME=/home/node/.openclaw/.local/share
Environment=XDG_CACHE_HOME=/home/node/.openclaw/.cache
Environment=CODEX_HOME=/home/node/.openclaw/.codex
Environment=CODEX_SQLITE_HOME=/home/node/.openclaw/.codex/sqlite
Environment=HOMEBREW_CACHE=/home/node/.openclaw/.cache/Homebrew
Volume=openquad-documents-state:/home/node/.openclaw:rw,Z
Volume=openquad-documents-linuxbrew:/home/linuxbrew/.linuxbrew:rw,Z
Volume=openquad-documents-npm:/home/node/.npm:rw,Z
Volume=openquad-documents-cache:/home/node/.cache:rw,Z
Volume=openquad-documents-logs:/tmp/openclaw:rw,Z
Exec=openclaw gateway --allow-unconfigured

[Service]
Restart=on-failure
RestartSec=5s
TimeoutStartSec=15min

[Install]
WantedBy=default.target
