# templates/communications-calendar/openquad.container
[Unit]
Description=OpenQuad Communications Agent

[Container]
ContainerName=openquad-comms
Image=ghcr.io/myos-dev/openquad-comms:latest
Pull=missing
Network=openquad.network

PublishPort=127.0.0.1:18789:18789
UserNS=keep-id
NoNewPrivileges=true
DropCapability=ALL
ReadOnly=true
Tmpfs=/tmp:size=512M,mode=1777
WorkingDir=/home/node
Environment=HOME=/home/node
Environment=OPENQUAD_TEMPLATE=communications-calendar
Environment=OPENQUAD_IMAGE_NAME=openquad-comms
Environment=OPENCLAW_STATE_DIR=/home/node/.openclaw
Environment=OPENCLAW_CONFIG_PATH=/home/node/.openclaw/openclaw.json
Environment=OPENQUAD_DEFAULT_CONFIG=/usr/share/openquad/templates/communications-calendar/openclaw.json5
Environment=OLLAMA_BASE_URL=http://ollama:11434
Environment=SEARXNG_BASE_URL=http://searxng:8080
Environment=NPM_CONFIG_CACHE=/home/node/.openclaw/.npm
Environment=XDG_CONFIG_HOME=/home/node/.openclaw/.config
Environment=XDG_DATA_HOME=/home/node/.openclaw/.local/share
Environment=XDG_CACHE_HOME=/home/node/.openclaw/.cache
Environment=CODEX_HOME=/home/node/.openclaw/.codex
Environment=CODEX_SQLITE_HOME=/home/node/.openclaw/.codex/sqlite
Environment=HOMEBREW_CACHE=/home/node/.openclaw/.cache/Homebrew
Volume=openquad-comms-state:/home/node/.openclaw:rw,Z
Volume=openquad-comms-linuxbrew:/home/linuxbrew/.linuxbrew:rw,Z
Volume=openquad-comms-npm:/home/node/.npm:rw,Z
Volume=openquad-comms-cache:/home/node/.cache:rw,Z
Volume=openquad-comms-logs:/tmp/openclaw:rw,Z
Exec=openclaw gateway --allow-unconfigured

[Service]
Restart=on-failure
RestartSec=5s
TimeoutStartSec=15min

[Install]
WantedBy=default.target
