# templates/records/openquad.container
[Unit]
Description=OpenQuad Records Agent

[Container]
ContainerName=openquad-records
Image=ghcr.io/myos-dev/openquad-records:latest
Pull=missing
Network=openquad.network

PublishPort=127.0.0.1:18789:18789
UserNS=keep-id
NoNewPrivileges=true
DropCapability=ALL
ReadOnly=true
Tmpfs=/tmp:size=512M,mode=1777
WorkingDir=/home/node
Environment=HOME=/home/node
Environment=OPENQUAD_TEMPLATE=records
Environment=OPENQUAD_IMAGE_NAME=openquad-records
Environment=OPENCLAW_STATE_DIR=/home/node/.openclaw
Environment=OPENCLAW_CONFIG_PATH=/home/node/.openclaw/openclaw.json
Environment=OPENQUAD_DEFAULT_CONFIG=/usr/share/openquad/templates/records/openclaw.json5
Environment=OLLAMA_BASE_URL=http://ollama:11434
Environment=VIC_API_BASE_URL=http://vic-api:8080
Environment=DATABASE_URL=postgresql://openquad-records@postgres/openquad_records
Environment=NPM_CONFIG_CACHE=/home/node/.openclaw/.npm
Environment=XDG_CONFIG_HOME=/home/node/.openclaw/.config
Environment=XDG_DATA_HOME=/home/node/.openclaw/.local/share
Environment=XDG_CACHE_HOME=/home/node/.openclaw/.cache
Environment=CODEX_HOME=/home/node/.openclaw/.codex
Environment=CODEX_SQLITE_HOME=/home/node/.openclaw/.codex/sqlite
Environment=HOMEBREW_CACHE=/home/node/.openclaw/.cache/Homebrew
Volume=openquad-records-state:/home/node/.openclaw:rw,Z
Volume=openquad-records-linuxbrew:/home/linuxbrew/.linuxbrew:rw,Z
Volume=openquad-records-npm:/home/node/.npm:rw,Z
Volume=openquad-records-cache:/home/node/.cache:rw,Z
Volume=openquad-records-logs:/tmp/openclaw:rw,Z
Exec=openclaw gateway --allow-unconfigured

[Service]
Restart=on-failure
RestartSec=5s
TimeoutStartSec=15min

[Install]
WantedBy=default.target
