# templates/browser-runtime-headless/openquad.container
[Unit]
Description=OpenQuad Headless Browser Runtime

[Container]
ContainerName=openquad-browser-runtime-headless
Image=ghcr.io/myos-dev/openquad-browser-runtime-headless:latest
Pull=missing
Network=openquad.network

# Do not publish the Playwright WebSocket publicly. It is privileged
# browser-control access and should be reachable only by trusted agents.
NoNewPrivileges=true
ReadOnly=false
ShmSize=1g
Environment=PLAYWRIGHT_WS_HOST=0.0.0.0
Environment=PLAYWRIGHT_WS_PORT=3000
Environment=PLAYWRIGHT_HEADLESS=true
Environment=BROWSER_DOWNLOAD_DIR=/home/pwuser/downloads
Environment=BROWSER_ARTIFACTS_DIR=/home/pwuser/artifacts
Volume=openquad-browser-headless-downloads:/home/pwuser/downloads:rw,Z
Volume=openquad-browser-headless-artifacts:/home/pwuser/artifacts:rw,Z

[Service]
Restart=on-failure
RestartSec=5s
TimeoutStartSec=5min

[Install]
WantedBy=default.target
